# mpdignore.path
#mpdignore.path
[Install]
WantedBy=paths.target

[Path] 
PathChanged=/library/music/mpd playlists/.mpdignore.m3u
#PathChanged=/var/lib/mpd/playlists/.mpdignore.m3u 
#PathChanged=/library/music/mpd playlists/.mpdignore.m3u
