# pulsevol.service
[Unit]
Description=Volume Control Service

[Service]
ExecStart=/usr/local/bin/pulsevol %i
